# roles/haproxy/files/haproxy_container.service
[Unit]
After=docker.service

[Service]
ExecStart=/etc/haproxy/start_haproxy.sh

[Install]
WantedBy=default.target